# airootfs/etc/systemd/system/reflector.service.d/archiso.conf
[Service]
ExecStart=
ExecStart=/usr/bin/reflector --country US --protocol https --latest 5 --age 2 --sort rate --save /etc/pacman.d/mirrorlist
