# airootfs/etc/systemd/system/systemd-resolved.service.d/disable.conf
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-resolved
Restart=always
