# airootfs/etc/systemd/system/systemd-networkd.service.d/disable.conf
[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd
Restart=always
